# iac/oci-lottery/dist/phenotype-oci-lottery.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.phenotype.oci-lottery</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/oci-lottery</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>RUST_LOG</key>
    <string>info</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/Users/REPLACE_ME/Library/Logs/phenotype-oci-lottery.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/REPLACE_ME/Library/Logs/phenotype-oci-lottery.err.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/REPLACE_ME</string>
</dict>
</plist>
